# pi-monitor.service
[Unit]
Description=Pi System Monitor
After=network.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/pi-monitor
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
